# emag-sync-scheduler.service
[Unit]
Description=eMAG Sync Scheduler
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/magflow/MagFlow
Environment=PATH=/home/magflow/anaconda3/envs/MagFlow/bin
Environment=PYTHONPATH=/home/magflow/MagFlow
ExecStart=/home/magflow/anaconda3/envs/MagFlow/bin/python3 /home/magflow/MagFlow/sync_scheduler.py start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=emag-sync-scheduler

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/magflow/MagFlow/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
