# sync-monitor.service
[Unit]
Description=eMAG Sync Monitor and Recovery Service
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/opt/magflow
Environment=PYTHONPATH=/opt/magflow
ExecStart=/usr/bin/python3 /opt/magflow/sync_monitor_recovery.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=256M
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/magflow/logs

[Install]
WantedBy=multi-user.target
